# spring-cloud-demo.ini
[program:eureka]
command=java -jar /usr/local/project/eureka-demo.jar --spring.profiles.active=test

[program:user-service-invoker]
command=java -jar /usr/local/project/user-service-invoker.jar --spring.profiles.active=test

[program:user-service-provider]
command=java -jar /usr/local/project/user-service-provider.jar --spring.profiles.active=test

[program:zuul]
command=java -jar /usr/local/project/zuul-demo.jar --spring.profiles.active=test